# instruments/gni/scripts/ensure_running.cron
#!/bin/csh -x

# Start / keep running gni_server.py (relay program from GNI instrument to control.
ps ax | grep gni_server | grep -qv grep
if ($status) then
  nohup /home/pi/gni/scripts/gni_server.py >& /tmp/gni_server.log &
endif
